# systemd/sing-grok-singbox.service
[Unit]
Description=Local sing-box core for Sing Grok
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sing-grok
Group=sing-grok
EnvironmentFile=/etc/sing-grok/backend.env
ExecStart=/opt/sing-grok/bin/sing-box run -c /etc/sing-grok/config.json
Restart=on-failure
RestartSec=3
LimitNOFILE=1048576
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/etc/sing-grok /var/lib/sing-grok /var/log/sing-grok
StandardOutput=append:/var/log/sing-grok/sing-box.log
StandardError=append:/var/log/sing-grok/sing-box.log

[Install]
WantedBy=multi-user.target
